# profile.d/me.sh
#!/usr/bin/env bash

export ME_ANSI_END="\e[0m"
export ME_ANSI_BOLD="\e[1m"
export ME_ANSI_RED="\e[91m"
export ME_ANSI_GREEN="\e[92m"
export ME_ANSI_YELLOW="\e[93m"
export ME_ANSI_BLUE="\e[94m"
export ME_ANSI_MAGENTA="\e[95m"
export ME_ANSI_CYAN="\e[96m"
export ME_ANSI_WHITE="\e[97m"


# environment variables
# -----------------------------------------------------------------------------
export EDITOR=vim
export VISUAL=vim
export FCEDIT=vim


# aliases
# -----------------------------------------------------------------------------
alias ls="ls --color=always"
alias ll="ls -lh"
