# etc/profile.d/me.sh
# environment variables
# -----------------------------------------------------------------------------
export EDITOR=vim
export VISUAL=vim
export FCEDIT=vim


# aliases
# -----------------------------------------------------------------------------
alias ls="ls --color=always"
alias ll="ls -lh"
