# etc/profile.d/me.sh
# environment variables
# -----------------------------------------------------------------------------
export EDITOR=vim
export VISUAL=vim
export FCEDIT=vim


# aliases
# -----------------------------------------------------------------------------
alias ls="ls --color=auto"
alias ll="ls -lh"
alias lsc="ls --color=always"
alias llc="lsc -lh"

